# hyops/drivers/config/ansible/collections/ansible_collections/hybridops/common/roles/decision_dispatcher/templates/decision-dispatcher.service.j2
[Unit]
Description=HybridOps Decision Dispatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ decision_dispatcher_run_as_user }}
Group={{ decision_dispatcher_run_as_group }}
WorkingDirectory={{ decision_dispatcher_root }}
ExecStart=/usr/bin/env python3 {{ decision_dispatcher_bin_path }}
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths={{ decision_dispatcher_root }}

[Install]
WantedBy=multi-user.target
