# hyops/drivers/config/ansible/collections/ansible_collections/hybridops/common/roles/edge_observability/templates/edge-observability.service.j2
[Unit]
Description=HybridOps Edge Observability (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ edge_obs_root }}
ExecStart=/usr/bin/docker compose -f {{ edge_obs_compose_path }} up -d
ExecReload=/usr/bin/docker compose -f {{ edge_obs_compose_path }} up -d
ExecStop=/usr/bin/docker compose -f {{ edge_obs_compose_path }} down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
