# hyops/drivers/config/ansible/collections/ansible_collections/hybridops/common/roles/decision_service/templates/decision-service.service.j2
[Unit]
Description=HybridOps Edge Decision Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ decision_run_as_user }}
Group={{ decision_run_as_group }}
WorkingDirectory={{ decision_root }}
ExecStart=/usr/bin/env python3 {{ decision_bin_path }}
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths={{ decision_root }}

[Install]
WantedBy=multi-user.target
