# hyops/drivers/config/ansible/collections/ansible_collections/hybridops/common/roles/decision_service/templates/decision-service.service.j2
[Unit]
Description=HybridOps Edge Decision Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ decision_service_run_as_user }}
Group={{ decision_service_run_as_group }}
WorkingDirectory={{ decision_service_root }}
Environment=HYOPS_RUNTIME_ROOT={{ decision_service_runtime_root }}
Environment=HYOPS_ENV={{ decision_service_runtime_env }}
EnvironmentFile=-{{ decision_service_env_path }}
ExecStart=/usr/bin/env python3 {{ decision_service_bin_path }}
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths={{ decision_service_root }}

[Install]
WantedBy=multi-user.target
