# hyops/drivers/config/ansible/collections/ansible_collections/hybridops/common/roles/decision_consumer/templates/decision-consumer.service.j2
[Unit]
Description=HybridOps Decision Consumer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 {{ decision_consumer_root }}/decision-consumer.py
Restart=always
RestartSec=5
WorkingDirectory={{ decision_consumer_root }}
User={{ decision_consumer_run_as_user }}
Group={{ decision_consumer_run_as_group }}

[Install]
WantedBy=multi-user.target
